# infrastructure/deployment/ansible/templates/etrid-validator.service.j2
[Unit]
Description=Ëtrid {{ chain_type | capitalize }} Validator ({{ validator_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ etrid_user | default('etrid') }}
Group={{ etrid_group | default('etrid') }}
WorkingDirectory=/opt/etrid

ExecStart=/opt/etrid/etrid \
  --name "{{ validator_name }}" \
  --chain {{ chain_name }} \
  --base-path /var/lib/etrid \
{% if 'validators' in group_names or 'backup_validators' in group_names %}
  --validator \
{% endif %}
  --port {{ p2p_port }} \
  --rpc-port {{ rpc_port }} \
  --ws-port {{ ws_port }} \
  --rpc-methods Unsafe \
  --rpc-cors all \
  --prometheus-external \
  --prometheus-port {{ substrate_metrics_port }} \
{% if telemetry_enabled %}
  --telemetry-url "{{ telemetry_url }} 0" \
{% endif %}
  --pruning {{ pruning }} \
  --db-cache {{ db_cache }} \
  --max-runtime-instances {{ max_runtime_instances }} \
{% if bootnodes | length > 0 %}
{% for bootnode in bootnodes %}
  --bootnode "{{ bootnode }}" \
{% endfor %}
{% endif %}
  --log info

# Restart configuration
Restart=always
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/etrid /var/log/etrid

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=etrid

[Install]
WantedBy=multi-user.target
