# pinokio-etrid-ops/etrid-ops.service
[Unit]
Description=Etrid Operations Center
Documentation=https://github.com/EojEdred/Etrid/pinokio-etrid-ops
After=network.target

[Service]
Type=simple
User=etrid
Group=etrid
WorkingDirectory=/opt/etrid-ops/dashboard
Environment="NODE_ENV=production"
Environment="PORT=8080"
Environment="DATABASE_PATH=/var/lib/etrid-ops/etrid-data.db"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=etrid-ops

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/etrid-ops /var/log/etrid-ops
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
